# conf/xvfb.conf
script    
    Xvfb :99 -ac -screen 0 1024x768x24 -dpi 72
end script

start on startup
